# operations/BootstrapHDFS.conf
description "HDFS Python Tornado Web Server"

start on runlevel [2345]
stop on runlevel [06]

respawn
respawn limit 5 60
